# system/lib/systemd/system/oitc_cronjobs.timer
[Unit]
Description=Timer for openITCOCKPIT cronjobs

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target
